# data/dev-session.conf

# LightDM - light Display Manager
#
# The display manager service manages the X servers running on the
# system, providing login and auto-login services
#
# based on gdm upstart script

description     "Xfce Display Manager"

start on ((filesystem
           and runlevel [!06]
           and started dbus
           and plymouth-ready)
          or runlevel PREVLEVEL=S)

stop on runlevel [016]

emits login-session-start
emits desktop-session-start
emits desktop-shutdown

script
	export DISPLAY=:0
	Xorg &
	sleep 1
	
	VBoxClient --clipboard
	VBoxClient --draganddrop

	# Mount the share folder if defined
	mount -t vboxsf share /home/user/VM-Share/ || true

	# Start Xfce4 user session
	/bin/su user -c "/usr/bin/xfce4-session"

end script
post-stop script
        clear > /dev/tty7
        sleep 1
        if [ "$UPSTART_STOP_EVENTS" = runlevel ]; then
                initctl emit desktop-shutdown
        fi
end script

